# roles/mosaico/templates/mosaico.service.j2
# {{ ansible_managed }}

[Unit]
Description=mosaico-standalone

[Service]
ExecStart=/usr/bin/npm start
WorkingDirectory={{ mosaico_repo }}
Restart=on-failure
Type=simple
StandardError=syslog
